# archive/perme_rewire_order36_20260718T222611Z_predeploy/plists/com.atlas.intraday.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.atlas.intraday</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/yasser/scripts/atlas_intraday_receipt_wrapper.py</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/yasser/scripts/atlas_intraday.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yasser/scripts/atlas_intraday.err.log</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Minute</key><integer>0</integer></dict>
        <dict><key>Minute</key><integer>10</integer></dict>
        <dict><key>Minute</key><integer>20</integer></dict>
        <dict><key>Minute</key><integer>30</integer></dict>
        <dict><key>Minute</key><integer>40</integer></dict>
        <dict><key>Minute</key><integer>50</integer></dict>
    </array>
</dict>
</plist>
