# archive/holdings_reunderwrite_production_20260717T123738Z/com.atlas.holdings_reunderwrite.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.atlas.holdings_reunderwrite</string>
  <key>ProgramArguments</key><array><string>/usr/bin/python3</string><string>/Users/yasser/scripts/atlas_holdings_reunderwrite_runner.py</string></array>
  <key>RunAtLoad</key><false/>
  <key>StartCalendarInterval</key><array>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>2</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>2</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>2</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>2</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Weekday</key><integer>7</integer><key>Hour</key><integer>2</integer><key>Minute</key><integer>20</integer></dict>
  </array>
  <key>ExitTimeOut</key><integer>260</integer>
  <key>StandardOutPath</key><string>/Users/yasser/Library/Logs/Atlas/holdings_reunderwrite.out.log</string>
  <key>StandardErrorPath</key><string>/Users/yasser/Library/Logs/Atlas/holdings_reunderwrite.err.log</string>
</dict></plist>
